# logstash/pipeline/windows-events.conf
input {
  beats {
    port => 5044
    type => "windows"
    tags => ["siem_detected_event"]
  }

  udp {
    port => 5514
    type => "windows_syslog"
    tags => ["siem_detected_event"]
  }
}

filter {
  if [type] == "windows" or [event][module] == "security" {
    mutate {
      add_tag => ["windows_event"]
      add_field => { "siem_source" => "windows" }
    }

    # Windows Event ID based detection
    if [winlog][event_id] {

      # 4625 = Failed Logon
      if [winlog][event_id] == 4625 {
        mutate {
          add_tag => ["windows_failed_logon"]
          add_field => {
            "mitre.tactic" => "Credential Access"
            "mitre.technique_id" => "T1110"
            "mitre.technique_name" => "Brute Force"
            "risk_score" => 30
          }
        }
      }

      # 4624 = Successful Logon
      if [winlog][event_id] == 4624 {
        mutate {
          add_tag => ["windows_logon_success"]
          add_field => { "risk_score" => 5 }
        }
      }

      # 4688 = Process Created
      if [winlog][event_id] == 4688 {
        mutate {
          add_tag => ["windows_process_create"]
          add_field => {
            "mitre.tactic" => "Execution"
            "mitre.technique_id" => "T1059"
            "risk_score" => 10
          }
        }
        # Check for PowerShell without /i flag (Logstash regex is case-sensitive)
        if [winlog][event_data][NewProcessName] =~ /[Pp]ower[Ss]hell\.exe/ {
          mutate {
            add_tag => ["powershell_exec"]
            add_field => { "risk_score" => 40 }
          }
        }
      }

      # 4698 = Scheduled Task Created
      if [winlog][event_id] == 4698 {
        mutate {
          add_tag => ["scheduled_task_created"]
          add_field => {
            "mitre.tactic" => "Persistence"
            "mitre.technique_id" => "T1053.005"
            "risk_score" => 50
          }
        }
      }

      # 4720 = User Account Created
      if [winlog][event_id] == 4720 {
        mutate {
          add_tag => ["user_account_created"]
          add_field => {
            "mitre.tactic" => "Persistence"
            "mitre.technique_id" => "T1136.001"
            "risk_score" => 60
          }
        }
      }

      # 7045 = New Service Installed
      if [winlog][event_id] == 7045 {
        mutate {
          add_tag => ["new_service_installed"]
          add_field => {
            "mitre.tactic" => "Persistence"
            "mitre.technique_id" => "T1543.003"
            "risk_score" => 55
          }
        }
      }
    }

    # GeoIP enrichment on Windows logon source IP
    if [winlog][event_data][IpAddress] and [winlog][event_data][IpAddress] != "-" {
      mutate {
        add_field => { "source_ip" => "%{[winlog][event_data][IpAddress]}" }
      }
      geoip {
        source => "source_ip"
        target => "geo"
      }
    }
  }
}

output {
  if "windows_event" in [tags] {
    elasticsearch {
      hosts => ["http://elasticsearch:9200"]
      index => "siem-logs-windows-%{+YYYY.MM.dd}"
      user => "elastic"
      password => "${ELASTIC_PASSWORD}"
    }
  }
}
